# etc/systemd/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin panjerd --noclear %I 115200 linux
Nice=-15
